# etc/systemd/system/openvpn@.service.d/suspend.conf
[Unit]
Conflicts=suspend.target sleep.target

[Service]
Restart=always
